# packaging/autostart/homeui.desktop
[Desktop Entry]
Type=Application
Name=HomeUI
Comment=Native OpenHAB touchscreen dashboard
Exec=/usr/local/bin/homeui --config /etc/homeui/dashboard.json
Terminal=false
X-GNOME-Autostart-enabled=true
NoDisplay=true
Categories=Utility;
